# clinic-bot.service
[Unit]
Description=Telegram Bot "Объединяем компетенции"
After=network.target

[Service]
Type=simple
User=clinicbot
Group=clinicbot
WorkingDirectory=/opt/clinic-bot
Environment="PATH=/opt/clinic-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/clinic-bot/venv/bin/python /opt/clinic-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clinic-bot

# Ограничения ресурсов
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
