# clinic-bot.service
[Unit]
Description=Telegram Bot "Объединяем компетенции"
After=network.target

[Service]
Type=simple
User=clinicbot
WorkingDirectory=/opt/clinic-bot
ExecStart=/usr/bin/python3 /opt/clinic-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clinic-bot

# Ограничения ресурсов
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
